# ports/inetutils/.footprint
drwxr-xr-x	0/0	bin/
-rwxr-xr-x	0/0	bin/dnsdomainname
-rwxr-xr-x	0/0	bin/domainname
-rwxr-xr-x	0/0	bin/hostname
drwxr-xr-x	0/0	etc/
drwxr-xr-x	0/0	etc/cron/
drwxr-xr-x	0/0	etc/cron/daily/
-rwxr-xr-x	0/0	etc/cron/daily/syslog
drwxr-xr-x	0/0	etc/rc.d/
-rwxr-xr-x	0/0	etc/rc.d/syslogd
-rw-r--r--	0/0	etc/syslog.conf
drwxr-xr-x	0/0	sbin/
-rwxr-xr-x	0/0	sbin/ifconfig
drwxr-xr-x	0/0	usr/
drwxr-xr-x	0/0	usr/bin/
-rwxr-xr-x	0/0	usr/bin/ftp
-rwxr-xr-x	0/0	usr/bin/telnet
-rwsr-xr-x	0/0	usr/bin/traceroute
-rwxr-xr-x	0/0	usr/bin/whois
drwxr-xr-x	0/0	usr/man/
drwxr-xr-x	0/0	usr/man/man1/
-rw-r--r--	0/0	usr/man/man1/dnsdomainname.1.gz
-rw-r--r--	0/0	usr/man/man1/domainname.1.gz
-rw-r--r--	0/0	usr/man/man1/ftp.1.gz
-rw-r--r--	0/0	usr/man/man1/hostname.1.gz
-rw-r--r--	0/0	usr/man/man1/ifconfig.1.gz
-rw-r--r--	0/0	usr/man/man1/telnet.1.gz
-rw-r--r--	0/0	usr/man/man1/traceroute.1.gz
-rw-r--r--	0/0	usr/man/man1/whois.1.gz
drwxr-xr-x	0/0	usr/man/man8/
-rw-r--r--	0/0	usr/man/man8/syslogd.8.gz
drwxr-xr-x	0/0	usr/sbin/
-rwxr-xr-x	0/0	usr/sbin/rotatelog
-rwxr-xr-x	0/0	usr/sbin/syslogd
